# netshaper.service
[Unit]
Description=Network Shaper Service (tc/netem + OLED display + HTTP API)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run as root so we can call `tc` and access /dev/i2c-1
User=root
ExecStart=/usr/local/bin/netshaper
Restart=on-failure
RestartSec=5

# Allow raw socket / network admin capability
AmbientCapabilities=CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
